# pbm-functional/pytest/conf/supervisord.d/mongod.ini
[program:mongod]
command=/usr/bin/mongod --bind_ip_all --dbpath /var/lib/mongo %(ENV_MONGODB_EXTRA_ARGS)s
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=mongod
