# pbm-functional/pytest/conf/supervisord.d/mongod-ce.ini
[program:mongod-ce]
command=/opt/mongodb/bin/mongod --wiredTigerCacheSizeGB 0.25 --bind_ip_all --dbpath %(ENV_DATADIR)s %(ENV_MONGODB_EXTRA_ARGS)s
autostart=%(ENV_AUTOSTART_CE)s
autorestart=false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=mongodb
